# crontab
*/1 * * * * /usr/local/bin/python /usr/src/app/insertrandomrecord.py 
