# crontab
*/1 * * * * python /insertrandomrecord.py